# pc/autostart/startimwheel.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/imwheel.sh
Hidden=false
Name=startimwheel
